# logstash/service/export_metrics/run
#!/bin/bash
set -e

exec /opt/logstash/bin/logstash -f /data/logstash/resources/export_metrics.conf
