# elasticsearch/etc/sv/kibana/run
#!/bin/sh

unset http_proxy && exec /home/kibana/bin/kibana > /dev/null 2>&1
