# elasticsearch/etc/sv/elasticsearch/run
#!/bin/sh
exec /home/elasticsearch/bin/elasticsearch > /dev/null 2>&1
