# logstash/service/logstash/finish
#!/bin/sh

kill $(pgrep logstash | head -1)