# elasticsearch/service/kibana/finish
#!/bin/sh

kill $(pgrep kibana | head -1)