# elasticsearch/service/kibana/run
#!/bin/sh

exec /home/kibana/bin/kibana > /dev/null 2>&1