# elasticsearch-base/service/elasticsearch/finish
#!/bin/sh

kill $(pgrep elasticsearch | head -1)